# dots/linux/.config/systemd/docker.service
[Unit]
Description=Start the Docker service

[Service]
Type=simple
ExecStart=/usr/local/bin/dockerd

[Install]
WantedBy=graphical.target

